# deploy/playbooks/roles/common/templates/circus.conf
start on filesystem and net-device-up IFACE=lo
setuid {{ ansible_ssh_user }}
exec {{ app_home }}/bin/circusd /etc/circus/circus.ini
